# systemd/user/app-sandbox-.service.d/dbus-session-bus.conf
[Unit]
Wants=xdg-dbus-proxy-filter@%p-%i.socket

[Service]
#RuntimeDirectory=%p/%i/dbus:dbus
BindReadOnlyPaths=%t/%p/%i/session-bus:%t/bus
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=%t/bus

