# systemd/user/app-sandbox-web-.service.d/10-networks.conf
[Service]
RestrictAddressFamilies=AF_INET AF_INET6
BindReadOnlyPaths=/run/systemd/resolve/stub-resolv.conf:/etc/resolv.conf
