# systemd/user/app-sandbox-web-.service.d/fido-tokens.conf
[Service]
# We need PrivateDevices=no in order to allow access to /dev/hidraw* for FIDO tokens.
PrivateDevices=no
DevicePolicy=closed
DeviceAllow=char-hidraw
RestrictAddressFamilies=AF_NETLINK
