# systemd/user/foot-server.service.d/graphical-slice.conf
[Unit]
After=graphical-session.target
[Service]
Slice=app-term.slice
