# systemd/user/app-sandbox-web-librewolf@.service
[Unit]
Description=Web Browser
CollectMode=inactive-or-failed

[Service]
ExecStart=librewolf
ExitType=cgroup
