# systemd/user/xdg-dbus-proxy@.socket
[Unit]
Description=D-Bus User Message Bus Socket

[Socket]
ListenStream=%t/%i/bus

[Install]
WantedBy=sockets.target
