# systemd/user/app-sandbox-meet-.service.d/pipewire.conf
[Service]
BindPaths=-%t/pulse/native
BindPaths=-%t/pipewire-0
BindPaths=-%t/pipewire-0.lock
