# systemd/user/app-sandbox-meet-.service.d/v4l.conf
[Service]
DeviceAllow=char-video4linux
DeviceAllow=char-media
