# systemd/user/app-sandbox-web-librewolf@.service.d/librewolf-profile-dir.conf
[Service]
BindPaths=%S/%p/%i/librewolf:%h/.librewolf
