# systemd/user/app-sandbox-web-librewolf@.service.d/librewolf-overrides-cfg.conf
[Service]
BindReadOnlyPaths=-%E/isol/librewolf.overrides.cfg:%h/.librewolf/librewolf.overrides.cfg
