# systemd/user/app-sandbox-web-librewolf@.service.d/ublock-import.conf
[Service]
BindReadOnlyPaths=-%E/isol/ublock-settings.txt:%h/ublock-settings.txt
