# systemd/user/app-sandbox-offline-.service.d/offline.conf
[Service]
PrivateNetwork=yes
Slice=app-sandbox-offline.slice
IPAddressDeny=any
IPAddressAllow=localhost
